# config/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
with-contenv
/usr/sbin/nginx -g "daemon off;"

