# config/s6-overlay/s6-rc.d/bridge-outlook/run
#!/command/execlineb -P
with-contenv
cd /app/bridge-outlook
s6-setuidgid nobody
/usr/bin/node index.js

